# delivery_qos/data/etc/cron.d/delivery_scan
0 2 * * * root /usr/bin/delivery_scan &
0 9 * * * root pkill delivery_scan
